# runit/x-ui/run
#!/bin/sh
exec 2>&1
echo 'Starting x-ui...'
sleep 5
exec chpst -e /usr/local/x-ui/ /usr/local/x-ui/x-ui
