# runit/x-ui/run
#!/bin/sh
exec 2>&1
echo -e "Starting x-ui..."
cd /usr/local/x-ui/ && exec /usr/local/x-ui/x-ui