# runit/x-ui/run
#!/bin/sh
exec 2>&1
echo 'Starting x-ui...'
cd /usr/local/x-ui/ && exec /usr/local/x-ui/x-ui
